# deploy/radio.conf
[program:audio-engine]
command=/usr/local/bin/audio_engine --config /etc/rakuraku/audio_engine.toml
directory=/var/lib/rakuraku
user=radio
autostart=true
autorestart=true
startretries=3
redirect_stderr=true
stdout_logfile=/var/log/rakuraku/audio-engine.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5

[program:radio-backend]
command=/usr/local/bin/radio-backend
directory=/var/lib/rakuraku
user=radio
environment=RADIO_CONFIG="/etc/rakuraku/radio-backend.toml",RADIO_LOG_LEVEL="info"
autostart=true
autorestart=true
startretries=3
redirect_stderr=true
stdout_logfile=/var/log/rakuraku/radio-backend.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5

[group:rakuraku]
programs=audio-engine,radio-backend
priority=999
